# templates/default/systemd/zookeeper.service.erb
[Unit]
Description=Zookeeper server

[Service]
EnvironmentFile=/etc/default/zookeeper
ExecStart=<%= @zkserver_bin %> start
ExecStop=<%= @zkserver_bin %> stop
User=<%= @zkuser %>
Type=forking
Restart=on-failure

[Install]
WantedBy=multi-user.target
