# templates/default/systemd/zookeeper.service.erb
[Unit]
Description=Zookeeper server

[Service]
EnvironmentFile=/etc/default/zookeeper
ExecStart=<%= @zkserver_bin %> start
ExecStop=<%= @zkserver_bin %> stop
User=<%= @zkuser %>
Type=forking
Restart=on-failure
LimitNOFILE=<%= @zk_nfiles %>
TimeoutSec=<%= @zk_timeout %>

[Install]
WantedBy=multi-user.target
