# templates/default/init.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides: zookeeper
# Short-Description: Zookeeper
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Required-Start:
# Required-Stop:
# Should-Start:
# Should-Stop:
# chkconfig: 2345 95 20
# description: Zookeeper
# processname: zookeeper
### END INIT INFO

[ -f /etc/init.d/functions ] && . /etc/init.d/functions

# TODO: Gotta do this more cleanly
<% node["zookeeper"]["env_vars"].each do |key, value| -%>
export <%= key %>=<%= value %>
<%- end -%>

# FIXME: Overall this init script needs some love
daemon --user=<%= node["zookeeper"]["user"] %> <%= node["zookeeper"]["base_directory"] %>/bin/zkServer.sh $1
