# templates/default/debian/systemd/zookeeper.service.erb
[Unit]
Description=Zookeeper server

[Service]
EnvironmentFile=- FIXME
ExecStart= FIXME
ExecReload= FIXME
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
