# templates/default/debian/systemd/zookeeper.service.erb
[Unit]
Description=Zookeeper server

[Service]
ExecStart=<%= @zkserver_bin %> start
ExecStop=<%= @zkserver_bin %> stop
ExecReload=<%= @zkserver_bin %> reload
User=<%= @zkuser %>
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
